# deploy/aggershield-server.service
[Unit]
Description=AggerShield control plane (licensing + admin dashboard)
Documentation=https://github.com/k4ran909/AggerShield
After=network-online.target
Wants=network-online.target

[Service]
# Put the admin token in /etc/aggershield/server.env:  AGGERSHIELD_ADMIN_TOKEN=...
EnvironmentFile=/etc/aggershield/server.env
ExecStart=/usr/local/bin/aggershield-server -listen :9000 -data /var/lib/aggershield/license.json
Restart=on-failure
RestartSec=2s
User=aggershield
Group=aggershield

AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/aggershield
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
